# debian/dist/matrix-os.service
[Unit]
Description=MATRIX OS Service
After=matrix-creator-malos.target

[Service]
User=matrix-os
Environment=MATRIX_MODE=service
ExecStart=/usr/lib/matrix-os/bin/node /usr/lib/matrix-os/index.js

[Install]
WantedBy=multi-user.target
