# debian/prerm
#!/bin/bash
echo "Disable firmware loading at startup"
systemctl disable matrix-os
