# debian/prerm
#!/bin/bash -e

# Stop, disable and remove matrix-os service
systemctl stop matrixio-os.service
systemctl disable matrixio-os.service


